# conf.d/chromium.conf
[program:chromium]
command=/usr/bin/chromium --no-sandbox --kiosk --disable-dev-shm-usage %(ENV_STARTING_WEBSITE_URL)s
autostart=true
autorestart=true