# conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -forever -usepw -shared -rfbport %(ENV_VNC_PORT)s -passwd "%(ENV_VNC_PASSWORD)s"
autostart=true
autorestart=true